# etc/inittab
::sysinit:/etc/rc.d/init.d/startup
::respawn:-/bin/sh

::shutdown:/etc/rc.d/shutdown
::ctrlaltdel:/bin/reboot
